# profile.d/00-aliases.sh
alias st='open -a "sublime text"'
alias sc='open -a "visual studio code"'

alias g='git'

alias v='vagrant'
alias vh='vagrant halt'
alias vs='vagrant suspend'
alias vu='vagrant up'
