# files/profile.d/00-aliases.sh
alias st='open -a "sublime text"'
alias sc='open -a "visual studio code"'

alias goo="open -a 'google chrome'"
alias ff="open -a 'firefox'"

alias g='git'
alias ya='yarn'

alias v='vagrant'
alias vh='vagrant halt'
alias vs='vagrant suspend'
alias vu='vagrant up'

alias ll='ls -alF'
alias ..='cd ../'

alias givkey="cat ~/.ssh/id_rsa.pub | tr -d '\n' | pbcopy"

alias python2="python2.7"
