# files/profile.d/03-git-completion-aliases.bash
# https://gist.github.com/JuggoPop/10706934

__git_complete g __git_main
__git_complete 'g co' _git_checkout
